# src/source/mkimage/features.in/gnome-atomic-customization/live/files/etc/profile.d/alias.sh
alias ll='ls -la'
alias ..='cd ..'
alias neofetch=fastfetch